# feeder/systemd/space-weather-mqtt.timer
[Unit]
Description=Space Weather MQTT timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target